# ops/systemd/saebooks-restore-test.service
[Unit]
Description=SAE Books weekly restore test
Documentation=file:///opt/saebooks/scripts/restore-test.sh
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/opt/saebooks/scripts/restore-test.sh
TimeoutStartSec=30min
StandardOutput=journal
StandardError=journal
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=default.target
